# config/systemd/debian-control/postinst
#!/bin/bash

systemctl enable aktualizr || true
systemctl start aktualizr || true

