# config/systemd/debian-control/prerm
#!/bin/bash

systemctl stop aktualizr || true
systemctl disable aktualizr || true

