# config/systemd/debian-control/prerm
#!/bin/bash

set -e

if [ "$1" = "remove" ]; then
    if [ -x "/usr/bin/deb-systemd-helper" ]; then
        deb-systemd-helper stop 'aktualizr.service' >/dev/null || true
    fi
fi
