# third_party/dbus-1.9.0/bus/dbus.socket
[Unit]
Description=D-Bus System Message Bus Socket

[Socket]
ListenStream=/usr/local/var/run/dbus/system_bus_socket
